# gggg.conf
[program:gggg]
command=python3 /home/curator2/gggg/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/gggg/test.err.log
stdout_logfile=/home/curator2/gggg/test.out.log